# debian/postinst
#!/bin/bash

set -e

systemctl daemon-reload
#systemctl enable sandbox
#systemctl start sandbox