# debian/prerm
#!/bin/bash

set -e

systemctl stop sandbox
#systemctl disable sandbox
